# .config/systemd/user/hour-telling.service
[Unit]
Description=Anunciar la hora actual

[Service]
Type=oneshot
ExecStart=/home/dilan/.local/bin/hourTelling
