# .config/systemd/user/clipmenud.service.d/override.conf
[Service]
ExecStartPost=/home/dilan/.local/bin/clipmenu-defaults.sh
