# etc/udev/rules.d/70-wol-realtek.rules
ACTION=="add", SUBSYSTEM=="net", ATTRS{idVendor}=="0bda", ATTRS{idProduct}=="8153", RUN+="/usr/bin/ethtool -s $name wol g"
