# etc/udev/rules.d/81-wol.rules
ACTION=="add", SUBSYSTEM=="net", NAME=="eth_monitor", RUN+="/usr/bin/ethtool -s eth_monitor wol g"
