# etc/udev/rules.d/99-intel-rapl.rules
SUBSYSTEM=="powercap", ACTION=="add", RUN+="/usr/bin/chmod -R 444 /sys/class/powercap/intel-rapl:0/"
